# stowed-files/i3/.config/autostart/polkit-agent.desktop
[Desktop Entry]
Type=Application
Name=PolicyKit Agent
Comment=
Exec=/usr/bin/lxqt-policykit-agent
OnlyShowIn=i3;
NoDisplay=true
